# packaging/syslog.service
[Unit]
Description=System Logging Service

[Service]
Type=forking
ExecStart=/sbin/syslogd

[Install]
WantedBy=multi-user.target
